# scripts/legacy/token_refresh.cron
# 42 Network Pipeline - Token Refresh
# Refresh API token every hour at minute 5

5 * * * * /srv/42_Network/repo/scripts/refresh_token.sh
